# hsz-origin-system-api/scripts/systemd/hsz-origin-nightly-check.timer
[Unit]
Description=Run HSZ Origin D-1/D-2 check at 04:30

[Timer]
OnCalendar=*-*-* 04:30:00 Asia/Shanghai
Persistent=true
Unit=hsz-origin-nightly-check.service

[Install]
WantedBy=timers.target
