# hsz-origin-system-api/scripts/systemd/hsz-origin-reconcile.service
[Unit]
Description=HSZ Origin nightly reconciliation
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hsz-origin/api
Environment=PYTHONPATH=/opt/hsz-origin/api
EnvironmentFile=/opt/hsz-origin/config/api.env
ExecStart=/opt/hsz-origin/api/venv/bin/python scripts/run_nightly_reconcile.py --max-repairs 24
StandardOutput=append:/opt/hsz-origin/logs/reconcile.log
StandardError=append:/opt/hsz-origin/logs/reconcile.log
